# daphne.service.example
[Unit]
Description=daphne daemon
After=network.target

[Service]
User=browny3542
Group=www-data
WorkingDirectory=/home/browny3542/Silk_Road
ExecStart=/home/browny3542/Silk_Road/venv/bin/daphne \
          -p 8001 \
          config.asgi:application

[Install]
WantedBy=multi-user.target

